# scripts/launchd/ai.reflectt.evi-presence-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.reflectt.evi-presence-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ryan/.openclaw/workspace/scripts/evi-presence-sync.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/reflectt-evi-presence-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/reflectt-evi-presence-sync-error.log</string>
</dict>
</plist>
